# res/init/systemd/stopwatchd.service
[Unit]
Description=Stopwatchd Daemon Server
Wants=

[Service]
ExecStart=<BINDIR>/swd
Type=simple
WorkingDirectory=/