# res/init/systemd/stopwatchd.service
[Unit]
Description=Stopwatchd Daemon Server
Wants=

[Service]
ExecStart=<PREFIX>/bin/swd
Type=simple
WorkingDirectory=/